# fs/home/pi/.config/autostart/touchegg.desktop
[Desktop Entry]
Type=Application
Exec=touchegg
